# debian/soundkid.service
[Unit]
Description=soundkid service

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/bin/soundkid

[Install]
WantedBy=multi-user.target
